# .config/superarch/systemd/theme.service
[Unit]
Description=Apply Superarch theme

[Service]
Type=oneshot
ExecStart=%h/.config/superarch/bin/switch-theme
